# monitor/scripts/rc.local
# Put your custom commands here that should be executed once
# the system init finished. By default this file does nothing.

./copytotmp

date=$(now)
echo $now & " ./copytotmp run on boot" >> /debug/logfile

exit 0
